# dotfiles/systemd/lightson@.service
[Unit]
Description=Lights On

[Service]
Environment=DISPLAY=%I
ExecStart=/bin/bash %h/.bin/lightson

[Install]
WantedBy=default.target
