# dotfiles/systemd/output-mail@.service
[Unit]
Description=output email for %i

[Service]
Type=oneshot
ExecStart=bash ${DOTFILES}/bin/systemd-output-mail %u@%H %i
