# dotfiles/systemd/status-mail@.service
[Unit]
Description=status email for %i

[Service]
Type=oneshot
ExecStart=bash ${DOTFILES}/../bin/systemd-mail %u@%H %i
