# dotfiles/systemd/weekly.timer
[Unit]
Description=Weekly timer

[Timer]
OnCalendar=weekly
Persistent=true
Unit=weekly.target

[Install]
WantedBy=timers.target
