# dotfiles/systemd/weekly.target
[Unit]
Description=Weekly timer target
StopWhenUnneeded=true
