# crontab/cron.daily/cleanup-trash
#!/bin/sh

command -v trash-empty >/dev/null 2>&1 || exit 0

trash-empty 30
