# crontab/cron.daily/logrotate
#!/bin/sh

/usr/sbin/logrotate --state $HOME/.var/lib/logrotate.status $HOME/.etc/logrotate.conf
