# crontab/cron.daily/logrotate
#!/bin/sh

[ -x /usr/sbin/logrotate ] || exit 0

/usr/sbin/logrotate --state "$HOME/.var/lib/logrotate.status" "$HOME/.etc/logrotate.conf"
