# crontab/cron.weekly/cleanup-vim-undo
#!/bin/bash

find $HOME/.var/lib/vim/undo -type f -mtime +100 -delete
