# systemd/daily.target
[Unit]
Description=Daily timer target
StopWhenUnneeded=true
