# systemd/geoclue-agent.service
[Unit]
Description=Start Geoclue agent

[Service]
ExecStart=/usr/lib/geoclue-2.0/demos/agent

[Install]
WantedBy=default.target
