# systemd/status-mail@.service
[Unit]
Description=status email for %i

[Service]
Type=oneshot
ExecStart=%h/.etc/bin/systemd-mail %u %i
