# crontab/run-parts
#!/bin/bash

set -u -o pipefail

if [[ -z "${1:-}" ]]; then
    echo "Usage: $0 <directory>" >&2
    exit 1
fi

for f in "$1"/*; do
    if [[ -x "$f" ]]; then
        "$f"
    fi
done
