# crontab/runcrontab
#!/bin/sh

PERIOD=$1

BASEDIR=$HOME/.etc/crontab
LOCALBASEDIR=$HOME/.local/etc/crontab

if command -v chrt >/dev/null && command -v ionice >/dev/null; then
    run() {
        chrt --idle 0 ionice -c2 -n7 "$@"
    }
else
    run() {
        nice "$@"
    }
fi

if command -v run-parts >/dev/null; then
    RUN_PARTS="run-parts --report"
else
    RUN_PARTS="$BASEDIR/run-parts"
fi

if test -d "$BASEDIR/cron.$PERIOD"; then
    run $RUN_PARTS "$BASEDIR/cron.$PERIOD"
fi

if test -d "$LOCALBASEDIR/cron.$PERIOD"; then
    run $RUN_PARTS "$LOCALBASEDIR/cron.$PERIOD"
fi
