# systemd/redshift-gtk.service.d/after.conf
# Also needs a [redshift] section in /etc/geoclue/geoclue.conf
[Unit]
After=geoclue-agent.service
